# systemd/greyline.timer
[Unit]
Description=Update the greyline world-time wallpaper on a schedule

[Timer]
OnCalendar=*:*:00
# Fire within 1s of :00 — a visible clock can't drift; one timer/min = negligible power.
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
